# mariadbp/rootfs/etc/services.d/ha-intergration/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: mariadb+ addon 
# Starts the home assisant intergration
# ==============================================================================

. /opt/mariadb-support/lib/mariadbio.sh
# ha Service Prefix from name...
SERVICEP=${HAA_NAME:-"mariadbp"}
# Secret for Service user for DB's 
#FIXME use VAR
SECF="/data/secret"
# PID FILE once ready
#FIXME use VAR
MDRF="/tmpfs/.mariadb_configure"
# dump location
DDATAL=${DDATAL:-"/data/dump"}
# backup location
BDBDIR=${BDBDIR:-"/backup"}

#Database
DATABASES=${DATABASES:-""}
RIGHTS=${RIGHTS:-""}
LOGINS=${LOGINS:-""}

WAITS=${WAITS:-120}
DUMPS=${DUMPS:-3}

CONFIG=/data/options.json

function ha-event() {
  local EVENTNAME="${SERVICEP}.${1}"
  local EVENTDATA=${2}

  bashio::log.info "Event ${EVENTNAME}"
  curl -sS -X "POST" \
    -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
    --data-binary "${EVENTDATA}" \
    "http://supervisor/core/api/events/${EVENTNAME}"
  echo
}

function stop() {
  if bashio::config.exists register;then
     for i in $(bashio::config register);do
       SERVICE=$i
       bashio::services.delete "${SERVICE}"
     done
  fi
}

function setup() {
  bashio::log.info "Waiting for mariadb..."
  while [[ ! -r "${MDRF}" ]];do
    # FIXME TIMEOUT
    sleep 1
  done
  wait_mariadb

  SERVICE="mysql"
  if bashio::services.available "${SERVICE}";then
      # Serice data cache
      SERVICEF="/tmpfs/${SERVICE}.service.json"
      DBSH=$(bashio::services "mysql" host)
      DBSP="$(bashio::services "mysql" port)"
      DBSU=$(bashio::services "mysql" username)
      DBSS=$(bashio::services "mysql" password)
      bashio::var.json \
        host "${DBSH}" \
        port "${DBSP}" \
        username "${DBSU}" \
        password "${DBSS}" > "${SERVICEF}"
      bashio::log.info "Service ${SERVICE} $(cat ${SERVICEF})"
      # my.cnf FIXME VAR remove on exit?
      cat <<EOF >"/config/mariadb/client.${SERVICE}.cnf"
[client]
host=${DBSH}
port=${DBSP}
user=${DBSU}
password=${DBSS}
EOF
      if ! bashio::fs.file_exists "${SECF}"; then
        # Clone password
         bashio::log.info "Cloning mysql service password"
        echo "${DBSS}" > "${SECF}"
      fi
  fi

  if ! bashio::fs.file_exists "${SECF}"; then
    bashio::log.info "Generating service password for Home Assistant"
    pwgen 64 1 > "${SECF}"
  fi
  #readonly 
  SECRET=$(<"${SECF}")
  #local PAYLOAD
  #readonly?
  PAYLOAD=$(bashio::var.json \
        host "$(hostname)" \
        port "3306" \
        username "service" \
        password "${SECRET}")

  #FXIME Private? Clone per service...
  bashio::log.info "Creating service user and permissions for Home Assistant"
  mysql -e "CREATE USER 'service'@'172.30.32.%' IDENTIFIED BY '${SECRET}';" 2> /dev/null || true
  mysql -e "SET PASSWORD FOR 'service'@'172.30.32.%' = PASSWORD('${SECRET}');" 2> /dev/null || true
  mysql -e "CREATE USER 'service'@'172.30.33.%' IDENTIFIED BY '${SECRET}';" 2> /dev/null || true
  mysql -e "SET PASSWORD FOR 'service'@'172.30.33.%' = PASSWORD('${SECRET}');" 2> /dev/null || true
  mysql -e "GRANT ALL PRIVILEGES ON *.* TO 'service'@'172.30.32.%' WITH GRANT OPTION;" 2> /dev/null || true
  mysql -e "GRANT ALL PRIVILEGES ON *.* TO 'service'@'172.30.33.%' WITH GRANT OPTION;" 2> /dev/null || true
  mysql -e "FLUSH PRIVILEGES;" 2> /dev/null || true

  #if mysql -e "SET PASSWORD FOR '${USERNAME}'@'%' = PASSWORD('${PASSWORD}');" 2> /dev/null; then
  #      echo "[INFO] Update user ${USERNAME}"
  #  else
  #      echo "[INFO] Create user ${USERNAME}"
  #      mysql -e "CREATE USER '${USERNAME}'@'%' IDENTIFIED BY '${PASSWORD}';" 2> /dev/null || true
  #fi

  # Init databases
  bashio::log.info "Init custom database(s)"
  for database in $DATABASES; do
     bashio::log.info "Create database ${database}"
    mysql -e "CREATE DATABASE ${database};" 2> /dev/null || true
  done

  # Init logins
  echo "[INFO] Init/Update users"
  for (( i=0; i < "$LOGINS"; i++ )); do
    USERNAME=$(jq --raw-output ".logins[$i].username" "$CONFIG")
    PASSWORD=$(jq --raw-output ".logins[$i].password" "$CONFIG")
    # Handel null password
    if [ "${PASSWORD}" = "null" ];then
       if [ ! -r "${SECF}.${USERNAME}" ];then
           pwgen 64 1 > "${SECF}.${USERNAME}"
       fi
       PASSWORD=$(<"${SECF}.${USERNAME}")
    fi

    if mysql -e "SET PASSWORD FOR '${USERNAME}'@'%' = PASSWORD('${PASSWORD}');" 2> /dev/null; then
         bashio::log.info "Update user ${USERNAME} with ${PASSWORD}"
    else
         bashio::log.info "Create user ${USERNAME} pwd ${PASSWORD}"
        mysql -e "CREATE USER '${USERNAME}'@'%' IDENTIFIED BY '${PASSWORD}';" 2> /dev/null || true
    fi
    #Write client config..
    if [ ! -r "/config/mariadb/client.${USERNAME}.cnf" ];then
      cat <<EOF >"/config/mariadb/client.${USERNAME}.cnf"
[client]
host=${DBSH}
port=${DBSP}
user=${USERNAME}
password=${PASSWORD}
EOF
    fi
  done

  # Init rights
  bashio::log.info "Init/Update rights"
  for (( i=0; i < "$RIGHTS"; i++ )); do
    USERNAME=$(jq --raw-output ".rights[$i].username" "$CONFIG")
    DATABASE=$(jq --raw-output ".rights[$i].database" "$CONFIG")
    # FIXME default grant for compat with std plugin
    GRANT=$(jq --raw-output ".rights[$i].grant" "$CONFIG")

    bashio::log.info  "Alter rights for ${USERNAME} to ${DATABASE}"
    mysql -e "GRANT ${GRANT} ${DATABASE}.* TO '${USERNAME}'@'%';" || true
  done
  mysqladmin shutdown

  wait_mariadb
  # Send service information to the Supervisor
  if bashio::config.exists register;then
     for i in $(bashio::config register);do
       SERVICE=$i
       bashio::log.debug "Service Publish ${SERVICE} DATA: ${PAYLOAD}"
       if bashio::services.publish "${SERVICE}" "${PAYLOAD}"; then
         bashio::log.info "Successfully send ${SERVICE} service information to Home Assistant."
         # FIXME Check registred? allways seems to suc even when failed
       else
         bashio::log.warning "Service ${SERVICE} message to Home Assistant failed!"
       fi
     done
  fi
  trap "stop" SIGTERM SIGINT
  trap "" SIGHUP

  bashio::net.wait_for 8123 homeassistant 120
  ha-event "service" "${PAYLOAD}"

  if [ "${DUMPS}" -gt 0 ];then
    ls -tp "${DDATAL}" | grep -v '/$' | tail -n +$((DUMPS+1))| while IFS= read -r f; do  bashio::log.info "purge ${DDATAL}/$f"; rm "${DDATAL}/$f"; done || true
  fi
  if [ "${DUMPS}" -gt 0 ];then
    ls -tAp "${DDATAL}" | grep -v '/$' | grep '^\.' | while IFS= read -r f; do  bashio::log.warn "purge ${DDATAL}/$f"; rm "${DDATAL}/$f"; done || true
  fi
}

function eventloop() {
  bashio::log.info 'Starting the Home Assistant STDIN service...'
  # shellcheck disable=SC2162
  while read msg; do
    bashio::log.debug "> ${msg}"
    msg="${msg%\"}"
    msg="${msg#\"}"
    cmd="$(echo "$msg" | jq -r '.command|values' || true)"
    bashio::log.info 'Received message command to execute:'
    bashio::log.info "> ${cmd}"
    case "${cmd}" in
      ping) ha-event "ping" "{ }" ;;
      restart) exit 0 ;;
      service) ha-event "service" "${PAYLOAD}" ;;
      shutdown)  # dup?
        ha-event "shutdown" "{ }"
        s6-svscanctl -t /var/run/s6/services
        exit 0 
      ;;
      dump)
        dump_mariadb "${DDATAL}"
        ha-event "dump" "{ }"
      ;;
      backup)
        #FIXME opt location
        dump_mariadb
        ha-event "backup" "{ }"
      ;;
      restore)
        restore_mariadb
        ha-event "restored" "{ }"
      ;;
      size)
        mysql -e 'SELECT table_schema AS "Database", ROUND(SUM(data_length + index_length) / 1024 / 1024, 2) AS "Size (MB)" FROM information_schema.TABLES GROUP BY table_schema;' | tee /tmpfs/mariadb_size.txt || true
        ha-event "size" '{ "size" : "'"${result}"'", "output" : "'"$(</tmpfs/mariadb_size.txt)"'" }'
      ;;
      check)
	#FIXME ARGS
        local result
        if mysqlcheck --all-databases | tee /tmpfs/mariadb_check_results; then
          bashio::log.info "Command completed successfully."
          result="true"
        else
          bashio::log.error "Check Command exited with a failure status."
          result="false"
        fi
        ha-event "checked" '{ "ok" : "'"${result}"'", "output" : "'"$(</tmpfs/mariadb_check_results.txt)"'" }'
      ;;
      migrate)
	# Get creds
	# do dump
	# load dump
      ;;
      copy)
        # arg db?
      ;;
      *) bashio::log.error "Command not found: ${cmd}" ;;
    esac
  done < /proc/1/fd/0
  bashio::log.info 'Finished the Home Assistant STDIN service...'
  exit 1
}

function main() {
	if [ ! -r /tmpfs/.mariadb_ready ]; then
	  setup
	fi
	eventloop
}

main
# vim: ft=sh
